# core/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin xodex --login-program /usr/local/bin/xodex-session --noclear %I $TERM
